# bpi/inkstand.timer
[Unit]
Description=Run inkstand per minute

[Timer]
AccuracySec=1s
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
