# bpi/inkstand.service
[Unit]
Description=Inkstand

[Service]
Type=simple
ExecStart=/root/inkstand/run_inkstand.sh
